# .config/systemd/user/backup-package-list.timer
[Unit]
Description=Run backup-package-list daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
